# build/flexget.conf
[program:flexget-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/flexget.sh
